# config/profile.d/easybuild-modules.csh
# add python startup script for interactive sessions
setenv MODULEPATH $MODULEPATH:/apps/easybuild/modules/setup

