# config/profile.d/easybuild-modules.csh
# add easybuild modules to module search
setenv MODULEPATH $MODULEPATH:/apps/modulefiles

