# config/profile.d/easybuild-modules.sh
# add easybuild modules to module search
export MODULEPATH=$MODULEPATH:/apps/easybuild/modules/setup
